# spider/spider/all_spider/crontab

0       3       *       *       1     /work/spider/all_spider/crontab_all_spider.sh
